# deploy/systemd/bot.service
[Unit]
Description=Animals Telegram Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/workspace/Animals
EnvironmentFile=/workspace/Animals/.env
Environment=PETS_ASSETS_ROOT=/workspace/Animals/assets/pets
ExecStart=/workspace/Animals/.venv/bin/python -m bot.telegram_main
Restart=on-failure
User=bot
Group=bot

[Install]
WantedBy=multi-user.target
