# deploy/systemd/bot.service
[Unit]
Description=Animals Telegram Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/workspace/Animals
EnvironmentFile=/workspace/Animals/.env
ExecStart=/usr/bin/python -m bot.main
Restart=on-failure
User=bot
Group=bot

[Install]
WantedBy=multi-user.target
